# code/supervisor_etc/supervisord/mytornado.ini
[program:mytornado]
command=python2.7 /home/USERNAME/webapps/mytornado/start_server.py --port=12345
